# files/etc/uci-defaults/99-disable-modem
#!/bin/sh
# Disable runtime cellular/modem services if any package accidentally exists.
# This does not release reserved memory by itself; reserved memory is handled by DTS patch at build time.

for svc in modemmanager rmtfs qrtr-ns qmi-modem-410-init DbusSmsForwardCPlus sms_tool; do
    if [ -x "/etc/init.d/$svc" ]; then
        "/etc/init.d/$svc" disable 2>/dev/null || true
        "/etc/init.d/$svc" stop 2>/dev/null || true
    fi
done

rm -f /etc/rc.d/S*modemmanager /etc/rc.d/S*rmtfs /etc/rc.d/S*qrtr* /etc/rc.d/S*qmi* /etc/rc.d/S*DbusSms* 2>/dev/null || true

exit 0
